# etc/udev/rules.d/monitor_hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0.0", ENV{XAUTHORITY}="/home/bronsa/.Xauthority", RUN+="/home/bronsa/bin/togglemonitor"
